# systemd-pkg/reeder-web.service
[Unit]
Description=Reeder web interface
After=network.target

[Service]
Type=simple
User=reeder
Group=reeder
ExecStart=/usr/lib/reeder/bin/reeder-web --config /etc/reeder/config.toml
WorkingDirectory=/var/lib/reeder
Environment=HOME=/var/lib/reeder
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=UV_PROJECT_ENVIRONMENT=/var/lib/reeder/.venv
Environment=UV_CACHE_DIR=/var/lib/reeder/.cache/uv
Environment=HF_HOME=/var/lib/reeder/.cache/huggingface
Environment=REEDER_WEB_HOST=0.0.0.0
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
